# files/system/usr/lib/systemd/system/system-flatpak-setup.timer
[Unit]
Description=Timer for system-flatpak-setup

[Timer]
OnBootSec=120

[Install]
WantedBy=timers.target